# src/slipp/generator/templates/roles/app-systemd-python/templates/systemd.service.j2
# Generated by slipp launch
[Unit]
Description={{ project_name }}-{{ service.name }}
After=network.target

[Service]
Type=simple
User={{ '{{ ansible_user }}' }}
Group={{ '{{ ansible_user }}' }}
WorkingDirectory=/opt/{{ project_name }}-{{ service.name }}
EnvironmentFile=-/opt/{{ project_name }}-{{ service.name }}/.env.production
Environment=PORT={{ service.port }}
ExecStart=/opt/{{ project_name }}-{{ service.name }}/.venv/bin/{{ execBinary }}{% if execScript %} {{ execScript }}{% endif %}{% if exec_args %} {{ exec_args }}{% endif %}

Restart=on-failure
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
